# airootfs/etc/systemd/system/set-theme.service
[Unit]
Description=set-theme-mat

[Service]
Type=oneshot
ExecStart=/etc/set-mat.sh

[Install]
WantedBy=multi-user.target
